# redis/etc/service/redis/run
#!/bin/sh
cd /redis-data
exec redis-server /etc/redis.conf